# packaging/systemd/coder-desktop-helper.service.in
[Unit]
Description=Coder Desktop VPN Helper
Documentation=https://coder.com/docs

[Service]
Type=dbus
BusName=com.coder.Desktop.Helper
ExecStart=@HELPER_BINARY_PATH@
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=dbus-com.coder.Desktop.Helper.service
